# deploy/systemd/unit3dwebup.service
[Unit]
Description=Unit3DWebUp FastAPI bot
Documentation=https://31december99.github.io/Unit3DWebUp/
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dev/Unit3DWebUp
Environment=DOCKER=false
Environment=PYTHONUNBUFFERED=1
Environment=REDIS_HOST=127.0.0.1
Environment=REDIS_PORT=6379
Environment=ENVPATH=%h/.config/unit3dprep
ExecStart=%h/dev/Unit3DWebUp/.venv/bin/uvicorn unit3dwup.start:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
